# root/etc/services.d/sockd/run
#!/usr/bin/with-contenv sh
exec sockd -f /etc/sockd.conf -N $WORKERS
